# recipes-support/recipes-ccsp/ccsp/ccsp-common-library/checkrpiwifisupport.path
[Unit]
Description=Check RPI Wifi Support Initialized Path

[Path]
PathExists=/tmp/pam_initialized
Unit=checkrpiwifisupport.service

[Install]
WantedBy=multi-user.target
